# stages/07.extra-tweaks/01.extra-scripts/01.eval-configuration/eval-config.service
[Unit]
Description=Run script for EVAL setup
After=network.target
Before=set-usb0-up.service

[Service]
Type=oneshot
ExecStart=/usr/bin/nmcli connection add type ethernet ifname usb0 con-name usb0
ExecStart=/usr/bin/nmcli con modify usb0 ipv4.method link-local
ExecStart=/usr/bin/adi-eval-config.sh
ExecStartPost=/bin/rm /etc/systemd/system/eval-config.service

[Install]
WantedBy=multi-user.target
