# stages/07.extra-tweaks/10.ad-r1m/files/slcan.target
[Unit]
Requires=slcan.service
After=slcan.service

